# deploy/systemd/ccc-agent.service
[Unit]
Description=CCC Agent Pull Service
After=network.target
Wants=network-online.target

[Service]
Type=oneshot
User=%i
EnvironmentFile=-%h/.ccc/.env
ExecStart=/bin/bash %h/.ccc/workspace/deploy/agent-pull.sh
StandardOutput=append:%h/.ccc/logs/pull.log
StandardError=append:%h/.ccc/logs/pull.log

[Install]
WantedBy=multi-user.target
