# deploy/systemd/ccc-dashboard.service
[Unit]
Description=CCC Dashboard (nginx serving WASM SPA)
Documentation=https://github.com/jordanhubbard/CCC
After=network.target ccc-server.service
Wants=ccc-server.service

[Service]
Type=forking
User=root
PIDFile=/run/nginx-ccc-dashboard.pid

# Serves the pre-built Leptos/WASM dist/ via nginx.
# Config: deploy/nginx/dashboard.conf  (proxies /api and /bus to ccc-server on 8789)
ExecStartPre=/usr/sbin/nginx -t -c /etc/nginx/sites-available/ccc-dashboard.conf
ExecStart=/usr/sbin/nginx -g 'daemon on;' -c /etc/nginx/sites-available/ccc-dashboard.conf
ExecReload=/bin/kill -s HUP $MAINPID
ExecStop=/bin/kill -s QUIT $MAINPID

Restart=on-failure
RestartSec=5s

StandardOutput=journal
StandardError=journal
SyslogIdentifier=ccc-dashboard

[Install]
WantedBy=multi-user.target
