# deploy/systemd/ccc-server.service
[Unit]
Description=CCC API Server (Rust/Axum)
Documentation=https://github.com/jordanhubbard/CCC
After=network.target

[Service]
Type=simple
User=jkh
WorkingDirectory=/home/jkh/Src/CCC

ExecStart=/usr/local/bin/ccc-server

# Core config
Environment=CCC_PORT=8789
Environment=CCC_DATA_DIR=/home/jkh/.ccc/data
Environment=RUST_LOG=ccc_server=info

# Load per-node secrets (QDRANT_FLEET_KEY, etc.)
EnvironmentFile=-AGENT_HOME/.ccc/.env

# Restart on crash, not clean exits
Restart=on-failure
RestartSec=5s

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ccc-server

[Install]
WantedBy=multi-user.target
