# deploy/systemd/heartbeat-local.service
[Unit]
Description=CCC Heartbeat Local Buffer (offline resilience)
After=network.target

[Service]
Type=simple
User=jkh
WorkingDirectory=/home/jkh/Src/CCC
ExecStart=/usr/bin/node deploy/heartbeat-local/server.mjs
Restart=always
RestartSec=5
Environment=HB_LOCAL_PORT=8790
EnvironmentFile=-/home/jkh/.ccc/.env

[Install]
WantedBy=multi-user.target
