# deploy/acc-repo-sync.service
[Unit]
Description=ACC repo sync — pull/commit/push shared AccFS repo
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=AGENT_USER
Group=AGENT_USER
WorkingDirectory=AGENT_HOME
ExecStart=/bin/bash %h/.acc/workspace/deploy/acc-repo-sync.sh
EnvironmentFile=-%h/.acc/.env
StandardOutput=append:%h/.acc/logs/repo-sync.log
StandardError=append:%h/.acc/logs/repo-sync.log

# Safety: don't let a stuck git command run forever
TimeoutStartSec=120
